# initramfs-mmc/etc/udev/rules.d/10-mmc.rules
# automount the mmc/sd card

ACTION=="add", KERNEL=="mmcblk[0-9]p[0-9]", RUN+="/etc/udev/scripts/mmc.sh %k"
ACTION=="remove", KERNEL=="mmcblk[0-9]p[0-9]", RUN+="/etc/udev/scripts/mmc.sh %k"
