# initramfs-mmc/etc/udev/rules.d/10-usb.rules
# automount the usb card

KERNEL=="sda[1-9]", RUN+="/etc/udev/scripts/usb.sh %k"
